# roles/etcd/templates/etcd.service.j2
[Unit] 
Description=etcd key-value store 
Documentation=https://github.com/etcd-io/etcd 
After=network-online.target 
Wants=network-online.target 

[Service] 
Type=notify 
WorkingDirectory={{ data_dir }}
EnvironmentFile={{ conf_dir }}/etcd.env 
ExecStart=/usr/local/bin/etcd 
Restart=always 
RestartSec=10s 
LimitNOFILE=40000 
User=etcd 
Group=etcd 
  
[Install] 
WantedBy=multi-user.target 

